# behavior/rules/macos/persistence_persistence_via_suspicious_launch_agent_or_launch_daemon.toml
[rule]
description = """
Identifies the creation of a launch agent or daemon property list file containing abnormal or suspicious values. An
adversary may establish persistence by installing a new launch agent or daemon which executes at login. This rule looks
within newly created property list files for a set of specific keys and values that contain certain file paths or system
binaries malware has commonly used when implementing persistence via this method.
"""
id = "c6037fad-ad13-46a6-9f7f-4deeef5ac69b"
license = "Elastic License v2"
name = "Persistence via Suspicious Launch Agent or Launch Daemon"
os_list = ["macos"]
reference = [
    "https://medium.com/red-teaming-with-a-blue-team-mentality/a-brief-look-at-macos-detections-and-post-infection-analysis-b0ede7ecfeb9",
    "https://objective-see.org/blog",
    "https://www.elastic.co/security-labs/DPRK-strikes-using-a-new-variant-of-rustbucket",
]
version = "1.0.53"

query = '''
persistence where event.action == "launch_daemon" and (Persistence.runatload == true or Persistence.keepalive == true) and
 process.executable != null and Persistence.args != null and
  Persistence.args like~ ("*/.*", 
                        "/private/tmp/*", 
                        "/private/var/root/Library/HTTPStorages/*",
                        "/private/var/root/Library/*",
                        "/var/tmp/*", 
                        "/var/lib/*",
                        "/tmp/*", 
                        "/var/folders/*", 
                        "/Users/Shared/*",
                        "/Users/*/Library/WebKit/*", 
                        "/Library/Containers/*", 
                        "/var/root/*", 
                        "/Library/WebServer/*",
                        "/Library/Graphics/*",
                        "/Library/Fonts/*",
                        "DYLD_INSERT_LIBRARIES*",
                        "LD_PRELOAD*",
                        "*/applet",
                        "*/Automator Application Stub",
                        "*/Application Stub",
                        "*python*", 
                        "*spawn*",
                        "*http*",
                        "sh*", 
                        "*/bin/sh", 
                        "bash*", 
                        "*/bin/bash", 
                        "zsh*", 
                        "*/bin/zsh", 
                        "*java*", 
                        "*curl*",
                        "*nscurl*",
                        "*osascript*",
                        "*nohup*",
                        "*wget*", 
                        "*exec*", 
                        "*base64*") and not
  Persistence.args like~ ("/Users/Shared/.Printix.net/Service/Printix Service.app/Contents/MacOS/Printix Service", 
                      "/var/tmp/LogMeIn Rescue - */Support-LogMeInRescue.app/Contents/MacOS/Support-LogMeInRescue", 
                      "/Users/*/.docker/run/docker.sock",
                      "/private/var/tmp/GoToAssist Remote Support Applet - 2EA1/GoToAssistRemoteSupport.app/Contents/MacOS/GoToAssistRemoteSupport",
                      "/Users/*/Applications/Grammarly Desktop.app//Contents/Library//LaunchAgents/Grammarly Desktop Helper.app/Contents/MacOS/Grammarly Desktop Helper",
                      "/Users/*/Library/Google/GoogleSoftwareUpdate/GoogleSoftwareUpdate.bundle/Contents/Resources/GoogleSoftwareUpdateAgent.app/Contents/MacOS/GoogleSoftwareUpdateAgent",
                      "/usr/local/jamf/bin/jamf", 
                      "/Users/*/Applications/.com.bomgar.bomgar-scc.*/Remote Support Customer Client.app/Contents/MacOS/bomgar-scc",
                      "/Applications/.com.bomgar.bomgar-scc.*/Remote Support Customer Client.app/Contents/MacOS/bomgar-scc",
                      "/Library/Addigy/go-agent",
                      "/usr/local/bin/elastic-agent",
                      "/Library/Bitdefender/AVP/product/bin/bdredline",
                      "https://elastic.jamfcloud.com", 
                      "/Users/*/Library/Application Support/*",
                      "/var/run/docker.sock",
                      "/opt/homebrew/opt/httpd/bin/httpd",
                      "/Library/.rpcnet/rpc.net", "/Library/.rpcnet/Ctes",
                      "/opt/homebrew/opt/jenkins-lts/libexec/jenkins.war",
                      "/Library/PrivilegedHelperTools/com.docker.socket",
                      "/opt/homebrew/opt/jenkins/bin/jenkins",
                      "/Library/airBalticBundle/.AnyConnect/fix_anyconnect_profile.sh",
                      "Location of the installed AutomationEngineSyncer executable (to be replaced by installer)", 
                      "/Applications/.com.bomgar.bomgar-scc/Remote Support Customer Client.app/Contents/MacOS/bomgar-scc",
                      "/opt/cisco/secureclient/bin/Cisco Secure Client - AnyConnect VPN Service.app/Contents/MacOS/vpnagentd",
                      "/Applications/Monotype Fonts/.Components/Services/CoreServices/MonotypeFontsMFEP",
                      "/Applications/Monotype Fonts/.Components/Services/CoreServices/MonotypeFontsMonitorService",
                        "rm -rf /usr/local/dgagent/; rm -rf /usr/local/lib/dgagent/; rm -f /Library/LaunchDaemons/com.verdasys.dgagent.plist /Library/LaunchAgents/com.verdasys.DgSessionSvc.plist;") and not
  Effective_process.executable like~ ("/Library/Application Support/JAMF/Jamf.app/Contents/MacOS/JamfDaemon.app/Contents/MacOS/JamfDaemon",
                                  "/usr/local/sbin/FileWave.app/Contents/MacOS/fwcld",
                                  "/usr/local/jamf/bin/jamf",
                                  "/Library/Intune/Microsoft Intune Agent.app/Contents/MacOS/IntuneMdmDaemon",
                                  "/Library/PrivilegedHelperTools/com.bomgar.bomgar-scc.ElevationHelper.daemon",
                                  "/Applications/MSP Anywhere Agent N-central.app/Contents/Resources/MSP Anywhere Helper",
                                  "/Library/Kandji/Kandji Agent.app/Contents/Helpers/Kandji Library Manager.app/Contents/MacOS/kandji-library-manager") and not
  process.executable like~ ("/System/*",
                        "/usr/local/sbin/FileWave.app/Contents/MacOS/fwcld",
                        "/Applications/Podman Desktop.app/Contents/MacOS/Podman Desktop",
                        "/Library/PrivilegedHelperTools/com.bomgar.bomgar-scc.ElevationHelper.daemon",
                        "/Library/PrivilegedHelperTools/*",
                        "/usr/local/munki/libexec/installhelper",
                        "/private/tmp/KSInstallAction.*/Install Google Software Update.app/Contents/Helpers/ksinstall") and not
  process.code_signature.signing_id like ("com.bomgar.bomgar-scc", 
                                       "com.apple.vim", 
                                       "com.apple.cat", 
                                       "com.apple.cfprefsd",
                                       "com.jetbrains.toolbox",
                                       "com.apple.pico",
                                       "ksinstall",
                                       "com.apple.shove",
                                       "GoToAssistRemoteSupportDaemon",
                                       "coreutils",
                                       "com.sublimetext.4",
                                       "com.apple.ditto",
                                       "NodewareAgent",
                                       "com.logmein.Rescue",
                                       "Antelope Launcher") and not
  Persistence.executable like "/Applications/AhnLab/ASTx/astxd" and 
  not process.executable : "/Applications/mcpproxy.app/Contents/MacOS/mcpproxy"
'''

min_endpoint_version = "8.11.0"
optional_actions = []
[[actions]]
action = "kill_process"
field = "process.entity_id"
state = 0

[[threat]]
framework = "MITRE ATT&CK"
[[threat.technique]]
id = "T1543"
name = "Create or Modify System Process"
reference = "https://attack.mitre.org/techniques/T1543/"
[[threat.technique.subtechnique]]
id = "T1543.001"
name = "Launch Agent"
reference = "https://attack.mitre.org/techniques/T1543/001/"

[[threat.technique.subtechnique]]
id = "T1543.004"
name = "Launch Daemon"
reference = "https://attack.mitre.org/techniques/T1543/004/"


[[threat.technique]]
id = "T1547"
name = "Boot or Logon Autostart Execution"
reference = "https://attack.mitre.org/techniques/T1547/"
[[threat.technique.subtechnique]]
id = "T1547.011"
name = "Plist Modification"
reference = "https://attack.mitre.org/techniques/T1547/011/"



[threat.tactic]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/tactics/TA0003/"

[internal]
min_endpoint_version = "8.11.0"
